# raspi/drone-control.service
[Unit]
Description=Drone WiFi control backend (MAVLink bridge + PWA server)
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/raspi
Environment=FC_PORT=/dev/serial0
Environment=FC_BAUD=57600
ExecStart=/home/pi/raspi/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
